# templates/etc/systemd/system/node-red.service.erb
[Unit]
Description=Node-RED service

[Service]
User=<%= scope['nodered::owner'] %>
Group=<%= scope['nodered::group'] %>
ExecStart=/usr/bin/node-red -u <%= scope['nodered::base_dir'] %> -p <%= scope['nodered::listen_port'] %>

[Install]
WantedBy=multi-user.target
Wants=postgresql.service
After=postgresql.service