# templates/etc/systemd/system/node-red.service.erb
[Unit]
Description=Node-RED service

[Service]
User=<%= @node-red_user %>
Group=<%= @node-red_user %>
ExecStart=/usr/bin/node-red -u <%= scope['node-red::base_dir'] %> -p <%= scope['node-red::listen_port'] %>

[Install]
WantedBy=multi-user.target
Wants=postgresql.service
After=postgresql.service